# guest/rootfs-overlay/etc/inittab
# gatOS guest inittab (busybox init — no openrc; see OS_PLAN.md T2.2).
# No getty: the serial console is a log file; access is SSH (or QGA for control).
# sshd: -D keep in foreground (busybox respawn owns the lifecycle), -e log to stderr
# (captured into the QEMU serial log). Config is the default /etc/ssh/sshd_config.
::sysinit:/sbin/init-gatos
::respawn:/usr/sbin/sshd -D -e
::respawn:/sbin/qga-gatos
::respawn:/sbin/sim-mount
::respawn:/sbin/mnt-mount
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/shutdown-gatos
